# django-app/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/src/django/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
